# pack/rootfs/etc/s6-overlay/s6-rc.d/gpustack/run
#!/command/with-contenv /bin/bash
# shellcheck shell=bash
# ================================
# GPUStack longrun service
# ================================

source /etc/profile
ARGS_FILE="/run/gpustack/args"

# Read arguments from the args file if it exists and is not empty
set --
if [ -s "$ARGS_FILE" ]; then
    while IFS= read -r line || [ -n "$line" ]; do
        [ -z "$line" ] && continue
        set -- "$@" "$line"
    done < "$ARGS_FILE"
fi

echo "[INFO] Starting gpustack."
if [ "$#" -gt 0 ]; then
    exec gpustack start "$@"
else
    exec gpustack start
fi
