# pack/s6-overlay/s6-rc.d/gpustack/finish
#!/command/with-contenv bash

echo "[ERROR] gpustack exited with code $1, shutting down all services..." >&2

s6-svscanctl -t /run/service
exit 1
