# pack/s6-overlay/s6-rc.d/gpustack-migration/up
/etc/s6-overlay/scripts/gpustack-migration.sh
