# pack/rootfs/etc/s6-overlay/s6-rc.d/grafana/finish
#!/command/with-contenv bash

SCRIPT_ROOT=/etc/s6-overlay/scripts
source "$SCRIPT_ROOT/base.sh"

EXIT_CODE="$1"
EXIT_SIGNAL="$2"

handleOptionalServiceExit "grafana" "$EXIT_CODE" "$EXIT_SIGNAL"
